# meta-firealarm/recipes-pingled/pingled/files/myservice.service
[Unit]


[Service]
Type=simple
ExecStart=/usr/bin/pingled.sh
Restart=always


[Install]
WantedBy=multi-user.target